# quadlet/worb-staging.container
# ==============================================================================
# worb CMS - Quadlet Configuration (Staging)
# Generated by CodeB CLI v3.3.0
# Date: 2025-12-22
# ==============================================================================
# Environment: staging
# Port: 4500
# Domain: worb-staging.codeb.kr
# ==============================================================================
# Infrastructure:
#   - App Server: app.codeb.kr (158.247.203.55)
#   - Database: db.codeb.kr (PostgreSQL 5432)
#   - Redis: db.codeb.kr (Redis 6379)
#   - WebSocket: ws.codeb.kr (Centrifugo 8000)
# ==============================================================================

[Unit]
Description=worb CMS Staging
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeblabdev-max/workb:staging
ContainerName=worb-staging
AutoUpdate=registry

# Network mode: host for direct server connections
Network=host

# Core environment
Environment=NODE_ENV=staging
Environment=PORT=4500
Environment=HOSTNAME=0.0.0.0

# Database (Storage Server: db.codeb.kr - separate staging database)
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb_staging?schema=public

# Redis (Storage Server: db.codeb.kr - separate staging prefix)
Environment=REDIS_URL=redis://db.codeb.kr:6379/1
Environment=REDIS_PREFIX=worb-staging:

# Centrifugo (Streaming Server: ws.codeb.kr)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# Secrets from external file
EnvironmentFile=/etc/containers/systemd/worb-staging.env

# Health check
HealthCmd=curl -sf http://localhost:4500/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits (reduced for staging)
PodmanArgs=--memory=1g --cpus=1

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
