# quadlet/worb.container
# ==============================================================================
# worb App - Quadlet Configuration (Production)
# Generated by CodeB CLI v3.0.3
# ==============================================================================
# Environment: production
# Port: 4000
# Network: host (Storage 서버 직접 연결)
# Database: db.codeb.kr (Storage 서버)
# WebSocket: ws.codeb.kr (Centrifugo)
# ==============================================================================

[Unit]
Description=worb CMS Production Container
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeb-kr/worb:latest
ContainerName=worb
AutoUpdate=registry

# Network: host 모드로 Storage 서버 직접 연결
Network=host

# 기본 환경 변수
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=HOSTNAME=0.0.0.0

# Storage 서버 연결 (db.codeb.kr - 64.176.226.119)
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb?schema=public
Environment=REDIS_URL=redis://db.codeb.kr:6379/0
Environment=REDIS_PREFIX=worb:

# Centrifugo 연결 (ws.codeb.kr - 141.164.42.213)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# 추가 환경 변수 (시크릿 포함) - 서버에 별도 생성 필요
EnvironmentFile=/etc/containers/systemd/worb.env

# Health Check
HealthCmd=curl -sf http://localhost:4000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# 리소스 제한
PodmanArgs=--memory=2g --cpus=2

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
