# quadlet/worb.container
# worb App - Quadlet Configuration
# Generated by CodeB CLI v3.0.3
# Environment: production
# Network: host (Storage 서버 직접 연결)
# Database: db.codeb.kr (Storage 서버)
# WebSocket: ws.codeb.kr (Centrifugo)

[Unit]
Description=worb App production Container
After=network-online.target

[Container]
Image=ghcr.io/${GITHUB_REPOSITORY_OWNER}/worb:latest
ContainerName=worb
Network=host
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=HOSTNAME=0.0.0.0

# Storage 서버 연결 (db.codeb.kr)
Environment=DATABASE_URL=postgresql://postgres:postgres@db.codeb.kr:5432/worb?schema=public
Environment=REDIS_URL=redis://db.codeb.kr:6379

# Centrifugo 연결 (ws.codeb.kr) - Socket.IO 대체
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api

# 추가 환경 변수 (시크릿 포함)
EnvironmentFile=/etc/containers/systemd/worb.env

HealthCmd=curl -sf http://localhost:4000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s
PodmanArgs=--memory=2g --cpus=2

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target