# quadlet/worb-staging.container
# worb App - Quadlet Configuration
# Generated by CodeB CLI v3.0.3
# Environment: staging
# Network: host (Storage 서버 직접 연결)
# Database: db.codeb.kr (Storage 서버, staging DB)
# WebSocket: ws.codeb.kr (Centrifugo)

[Unit]
Description=worb App staging Container
After=network-online.target

[Container]
Image=ghcr.io/${GITHUB_REPOSITORY_OWNER}/worb:staging
ContainerName=worb-staging
Network=host
Environment=NODE_ENV=staging
Environment=PORT=3002
Environment=HOSTNAME=0.0.0.0

# Storage 서버 연결 (db.codeb.kr) - staging DB
Environment=DATABASE_URL=postgresql://postgres:postgres@db.codeb.kr:5432/worb_staging?schema=public
Environment=REDIS_URL=redis://db.codeb.kr:6379/1

# Centrifugo 연결 (ws.codeb.kr) - Socket.IO 대체
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api

# 추가 환경 변수 (시크릿 포함)
EnvironmentFile=/etc/containers/systemd/worb-staging.env

HealthCmd=curl -sf http://localhost:3002/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s
PodmanArgs=--memory=1g --cpus=1

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target