# quadlet/worb-redis.container
# worb Redis - Quadlet Configuration
# Generated by CodeB CLI v7.0.8
# Environment: production
# Network: codeb-network (DNS-based container communication)

[Unit]
Description=worb Redis production Container
After=network-online.target

[Container]
Image=docker.io/library/redis:7-alpine
ContainerName=worb-redis
PublishPort=6401:6379
Network=codeb-network

Volume=worb-redis-data:/data:Z
HealthCmd=redis-cli ping
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target