# quadlet/worb-staging-blue.container
# ==============================================================================
# worb CMS - Quadlet Configuration (Staging Blue Slot)
# Generated by CodeB CLI v7.0.10
# Date: 2026-01-11
# ==============================================================================
# Environment: staging
# Slot: blue
# Port: 4500
# Domain: worb-staging.codeb.kr (when active)
# Preview: worb-staging-blue.codeb.kr (when inactive)
# ==============================================================================
# Blue-Green Deployment Strategy:
#   - Blue slot (this): Port 4500
#   - Green slot: Port 4501
#   - Traffic switching is managed by Caddy reverse proxy
#   - MCP API handles slot promotion and rollback
# ==============================================================================

[Unit]
Description=worb CMS Staging (Blue Slot)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeb-dev-run/worb:staging
ContainerName=worb-staging-blue
AutoUpdate=registry

# Network: codeb-network for container communication
Network=codeb-network
PublishPort=4500:3000

# Core environment
Environment=NODE_ENV=staging
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=SLOT=blue

# NextAuth (CRITICAL: Must match staging domain)
Environment=NEXTAUTH_URL=https://worb-staging.codeb.kr

# Database (Storage Server: db.codeb.kr - separate staging database)
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb_staging?schema=public

# Redis (Storage Server: db.codeb.kr - separate staging database index)
Environment=REDIS_URL=redis://db.codeb.kr:6379/1
Environment=REDIS_PREFIX=worb-staging:

# Centrifugo (Streaming Server: ws.codeb.kr)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# Secrets from external file
EnvironmentFile=/etc/containers/systemd/worb-staging.env

# Health check
HealthCmd=curl -sf http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits (reduced for staging)
PodmanArgs=--memory=1g --cpus=1

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
