# quadlet/worb-green.container
# ==============================================================================
# worb CMS - Quadlet Configuration (Green)
# Generated by CodeB CLI v4.0.0
# Date: 2026-01-05
# ==============================================================================
# Environment: production (Green slot)
# Port: 4001
# Domain: workb.net (via Caddy)
# ==============================================================================
# Blue-Green Deployment:
#   - Blue (4000): Primary production slot
#   - Green (4001): Secondary production slot
#   - Caddy switches between Blue/Green for zero-downtime
# ==============================================================================

[Unit]
Description=worb CMS Production (Green)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/codeb-dev-run/worb:latest
ContainerName=worb-green
AutoUpdate=registry

# Network: codeb-network for container communication
Network=codeb-network
PublishPort=4001:3000

# Core environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=DEPLOYMENT_SLOT=green

# NextAuth (CRITICAL: Must match production domain)
Environment=NEXTAUTH_URL=https://workb.net

# Database (Storage Server: db.codeb.kr)
Environment=DATABASE_URL=postgresql://worb_user:worb_secure_pass@db.codeb.kr:5432/worb?schema=public

# Redis (Storage Server: db.codeb.kr)
Environment=REDIS_URL=redis://db.codeb.kr:6379/0
Environment=REDIS_PREFIX=worb:

# Centrifugo (Streaming Server: ws.codeb.kr)
Environment=NEXT_PUBLIC_CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_URL=wss://ws.codeb.kr/connection/websocket
Environment=CENTRIFUGO_API_URL=http://ws.codeb.kr:8000/api
Environment=CENTRIFUGO_API_KEY=pRMupNs6HlGp7G6xkPsAFrI8hN4g6U0G
Environment=CENTRIFUGO_SECRET=of0KuRFjjzhq5LlBURCuKqzTUAA08hwL

# Secrets from external file
EnvironmentFile=/etc/containers/systemd/worb.env

# Health check
HealthCmd=curl -sf http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# Resource limits
PodmanArgs=--memory=2g --cpus=2

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target
